# deploy/systemd/beacon-worker@.service
[Unit]
Description=Beacon Email Worker (%i)
Documentation=https://github.com/your-org/beacon
After=network-online.target
Wants=network-online.target

# Usage:
#   systemctl enable --now beacon-worker@email-sendgrid.service
#   systemctl enable --now beacon-worker@email-mailgun.service
#
# %i expands to the instance name (the part after @), formatted as
# "<channel>-<provider>" (e.g. "email-sendgrid"). It is passed through as
# WORKER_SPEC and parsed into channel + provider so each worker instance
# routes through the correct SMTP config.

[Service]
Type=simple
# Shared Temporal / Infisical configuration for all worker instances.
EnvironmentFile=/etc/beacon/worker.env
# Override WORKER_SPEC with the instance name (%i).
Environment=WORKER_SPEC=%i
ExecStart=/usr/local/bin/email_worker
Restart=always
RestartSec=5s
# Run as a dedicated non-root user.
User=beacon
Group=beacon
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
