# deploy/systemd/beacon-worker@.service
[Unit]
Description=Beacon Email Worker (%i)
Documentation=https://github.com/your-org/beacon
After=network-online.target
Wants=network-online.target

# Usage:
#   systemctl enable --now beacon-worker@sendgrid.service
#   systemctl enable --now beacon-worker@mailgun.service
#
# %i expands to the instance name (the part after @), which is used as the
# PROVIDER_NAME so each worker instance routes through the correct SMTP config.

[Service]
Type=simple
# Shared Temporal / Infisical configuration for all worker instances.
EnvironmentFile=/etc/beacon/worker.env
# Override PROVIDER_NAME with the instance name (%i).
Environment=PROVIDER_NAME=%i
ExecStart=/usr/local/bin/email_worker
Restart=always
RestartSec=5s
# Run as a dedicated non-root user.
User=beacon
Group=beacon
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
